# .config/systemd/user/dock-watcher.service
[Service]
ExecStart=/bin/sh -c ${HOME}/bin/dock-watcher
Restart=always
NoNewPrivileges=false

[Install]
WantedBy=default.target
